# roles/os_configs/configure_wireguard_vpn/templates/wireguard_watcher.service
[Unit]
Description=WireGuard Watcher
After=network.target

[Service]
Type=simple
ExecStart=/etc/wireguard/wireguard-watcher.sh
Restart=always

[Install]
WantedBy=multi-user.target
